# service/mgmNode.conf
# This upstart service is the MOSES Master node

description "MOSES Slave Node"

start on runlevel [2345]
stop on runlevel [^2345]

exec su -c '/opt/mgmNode/mgmNode.py' - moses
